# deploy/systemd/timemanager-web.service
[Unit]
Description=TimeManager Pro Web
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/time-manager-pro
EnvironmentFile=/opt/time-manager-pro/.env
ExecStart=/opt/time-manager-pro/.venv/bin/gunicorn app.main:app -w 1 -k uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000 --log-level info --timeout 120 --graceful-timeout 30 --keep-alive 5 --access-logfile -
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target